# bin/cron/ditto-backup.cron
* * * * * root /opt/base2/ditto/bin/ditto
